# deploy/systemd/blinko-folder-drop.service
[Unit]
Description=Blinko Folder Drop Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blinko-drop
WorkingDirectory=/
ExecStart=/usr/local/bin/blinko-folder-drop run --config /etc/blinko-folder-drop/config.yaml
EnvironmentFile=-/etc/default/blinko-folder-drop
Restart=always
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
